# crontab
@reboot python3 ~/code/edu/python_udemy/block_fb
